# supervisord-openssh-server.conf
[program:openssh-server]
command=/usr/sbin/sshd -D
numprocs=1
autostart=true
autorestart=true
